# system_configuration/systemd/mivia-rover-platform.service
################################################################################
#
#  MIVIA Rover - Platform Bringup Systemd Service Unit
#  File: mivia-rover-platform.service
#  Version: 1.0
#  Last Updated: January 2026
#
################################################################################
#
#  DESCRIPTION
#  -----------
#  Systemd service unit for MIVIA Rover platform bringup. Defines service
#  startup behavior, dependencies, environment, and recovery policies.
#
#  FUNCTIONALITY
#  ----------
#  • Launches /etc/mivia_rover/scripts/start_mivia_rover.sh
#  • Runs as the detected user (not root) for security
#  • Sets working directory to ROS workspace root
#  • Loads environment configuration from /etc/mivia_rover/env/mivia_rover.env
#  • Depends on set-network.service (network pre-configuration)
#  • Auto-restarts on failure with 2-second delay
#  • Graceful shutdown with 15-second timeout (SIGINT signal)
#  • Enabled for multi-user.target (standard system boot)
#
#  LIFECYCLE
#  --------
#  Boot:      set-network.service → mivia-rover-platform.service
#  Shutdown:  Graceful SIGINT → 15s timeout → SIGKILL
#  Failure:   Wait 2s → Restart (infinite retries by default)
#
#  MANAGEMENT
#  ----------
#  # Check status
#  systemctl status mivia-rover-platform.service
#
#  # View logs
#  journalctl -u mivia-rover-platform.service -f
#
#  # Manual control
#  sudo systemctl restart mivia-rover-platform.service
#  sudo systemctl stop mivia-rover-platform.service
#  sudo systemctl start mivia-rover-platform.service
#
#  # Enable/disable autostart
#  sudo systemctl enable mivia-rover-platform.service
#  sudo systemctl disable mivia-rover-platform.service
#
#  DEPENDENCIES
#  -----------
#  • systemd (service manager)
#  • set-network.service (CAN network configuration)
#  • network.target and network-online.target
#  • /etc/mivia_rover/scripts/start_mivia_rover.sh
#  • /etc/mivia_rover/env/mivia_rover.env
#
#  INSTALLATION
#  -----------
#  Installed by: reload_services.sh
#  Target location: /etc/systemd/system/mivia-rover-platform.service
#
################################################################################

[Unit]
Description=MIVIA Rover Platform Bringup
Documentation=man:systemd.service(5)
Wants=network-online.target NetworkManager.service
After=network-online.target NetworkManager.service

# Ensure CAN/network pre-configuration completes first
Wants=set-network.service
After=set-network.service

[Service]
Type=simple

# Load runtime configuration (recommended to keep enabled)
EnvironmentFile=/etc/mivia_rover/env/mivia_rover.env

# Run as the detected user (filled by installer)
User=@MIVIA_ROVER_USER@

# Ensure a valid HOME for ROS 2
Environment=HOME=@MIVIA_ROVER_HOME@

# Start from workspace root (filled by installer)
WorkingDirectory=@MIVIA_ROVER_WS_PATH@

ExecStart=/etc/mivia_rover/scripts/start_mivia_rover.sh

Restart=on-failure
RestartSec=2s

TimeoutStopSec=15s
KillSignal=SIGINT

# (Optional but recommended) Avoid systemd killing children too early
KillMode=mixed

[Install]
WantedBy=multi-user.target